# systemd/tsbot-pulseaudio.service
[Unit]
Description=TSBot – PulseAudio Null-Sink (nur Audio; ts3client laeuft on-demand)
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=tsbot
Group=tsbot
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Umgebungsvariablen aus config.env laden
EnvironmentFile=/opt/tsbot/config/config.env
ExecStartPre=+/bin/bash -c 'mkdir -p /run/user/1000 && chown tsbot:tsbot /run/user/1000 && chmod 700 /run/user/1000'
# Hinweis: Skript wird nach /usr/local/bin/start_tsbot_pulseaudio.sh installiert (siehe README).
# Race-fest (Retry-Schleife). Startet NUR PulseAudio + Null-Sink.
# Xvfb und ts3client laufen on-demand via tsbot-ts3client.service/.path (NICHT mehr hier).
ExecStart=/usr/local/bin/start_tsbot_pulseaudio.sh
ExecStop=/bin/bash -c 'pulseaudio --kill 2>/dev/null; rm -f /run/user/1000/pulse/native 2>/dev/null; exit 0'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tsbot-pulseaudio

[Install]
WantedBy=multi-user.target
